# Fabric/node_exporter/upstart.conf
start on startup

script
   /usr/bin/node_exporter
end script
